# scripts/launchd/com.tbp.chico.consolidate-memory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tbp.chico.consolidate-memory</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/bissuh/Documents/TBP/chico/scripts/consolidate-memory.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bissuh/Documents/TBP/chico</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bissuh/Documents/TBP/chico/logs/launchd-consolidate.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/bissuh/Documents/TBP/chico/logs/launchd-consolidate.err</string>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
